# deploy/mycelium-seed.service
[Unit]
Description=Mycelium Network sporocarp (relay + DNS Spore Bank)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mycelium
Group=mycelium
EnvironmentFile=/etc/mycelium/seed.env
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/mycelium daemon \
  --listen /ip4/0.0.0.0/tcp/4001 \
  --listen /ip6/::/tcp/4001 \
  --announce-ip ${MYCELIUM_ANNOUNCE_IP} \
  --announce-ip6 ${MYCELIUM_ANNOUNCE_IP6} \
  --no-mdns \
  --sporocarp \
  --horizon-port 7474 \
  --contribute 2cpu,4gb,100gb
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
# Horizon fica em 127.0.0.1 — não exponha 7474 sem reverse proxy.
# DuckDNS (opcional): DUCKDNS_TOKEN + DUCKDNS_DOMAIN no seed.env

[Install]
WantedBy=multi-user.target
